# build/gerrit/root/etc/inittab
#
::once:sleep 300
::once:/opt/remove-progs.sh