# build/gerrit/root/etc/inittab
#
::sysinit:/opt/start-gerrit.sh
::once:java -jar /opt/gerrit.war daemon -d /data --console-log
